# overlay/etc/rc.local
#!/bin/sh
mount_nfs -onfsv4 [2a01:e140:10:10::a]:tests /.kyua
kyua test -k /usr/tests/Kyuafile
init 0
