# overlay/etc/rc.local
#!/bin/sh
cd /usr/tests
kyua test --results-file=/.kyua/results.db
kyua report-html --output=/.kyua/results --results-file=/.kyua/results.db
init 0
